# files/canvos/overlay/files/etc/systemd/system/stylus-agent.service.d/bcm-sync.conf
[Unit]
# Wait for BCM cmd daemon to set hostname before syncing to Palette
After=cmd.service

[Service]
# Pre-registration cleanup — no-op when already registered. Runs first.
ExecStartPre=/usr/bin/palette-cleanup-stale.sh
ExecStartPre=/usr/bin/bcm-sync-userdata.sh
